# infra/systemd/airtable-mcp-server.service
[Unit]
Description=Airtable MCP Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/mcp-servers/airtable-mcp-server
ExecStart=/usr/bin/node --experimental-modules dist/index.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=AIRTABLE_API_KEY=pat0dxG1HQcAUz2Kr.2c6134ed2b8df1435d7966f9434da6755d1467fd4a62fee0de7984b48f448ac9

[Install]
WantedBy=multi-user.target
