# infra/systemd/cloudflared-n8n-oauth2.service
[Unit]
Description=Cloudflare Tunnel for n8n OAuth2 Callback
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/cloudflared tunnel --url http://localhost:5678 --hostname n8n-oauth2.rensto.com
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
