# infra/systemd/cloudflared.service
# Cloudflare Tunnel systemd service
# Place this file at /etc/systemd/system/cloudflared.service

[Unit]
Description=Cloudflare Tunnel for SuperSeller AI
After=network.target
Wants=network-online.target
Before=docker.service

[Service]
Type=notify
TimeoutStartSec=0
Restart=on-failure
RestartSec=5s
ExecStart=/usr/local/bin/cloudflared tunnel --config /etc/cloudflared/config.yml run
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudflared
KillMode=mixed
KillSignal=SIGTERM

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/cloudflared /var/log/cloudflared
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Resource limits
LimitNOFILE=65535
LimitNPROC=512
CPUWeight=50
MemoryMax=512M
TasksMax=100

# Environment variables (optional)
Environment="TUNNEL_METRICS=0.0.0.0:2000"
Environment="TUNNEL_LOGLEVEL=info"

[Install]
WantedBy=multi-user.target

# Installation instructions:
# 1. Copy this file to /etc/systemd/system/cloudflared.service
# 2. Reload systemd: sudo systemctl daemon-reload
# 3. Enable service: sudo systemctl enable cloudflared
# 4. Start service: sudo systemctl start cloudflared
# 5. Check status: sudo systemctl status cloudflared
# 6. View logs: sudo journalctl -u cloudflared -f

# Useful commands:
# Restart: sudo systemctl restart cloudflared
# Stop: sudo systemctl stop cloudflared
# Disable: sudo systemctl disable cloudflared
# Edit: sudo systemctl edit cloudflared
