# infra/systemd/github-mcp-server.service
[Unit]
Description=GitHub MCP Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/mcp-servers/github-mcp-server
ExecStart=/usr/bin/python3 server.py
Restart=always
RestartSec=10
Environment=PYTHONPATH=/opt/mcp-servers/github-mcp-server
Environment=GITHUB_TOKEN=your-github-token

[Install]
WantedBy=multi-user.target
